# config/raspi-ai.service
[Unit]
Description=raspi-ai python bot
After=network-online.target sound.target

[Service]
Type=simple
User=pi
Group=pi
SupplementaryGroups=tty
WorkingDirectory=/opt/raspi-ai/app
ExecStart=/opt/raspi-ai/app/run.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=RASPI_AI_MEDIA_TEST=1
Environment=DISPLAY=:0
Environment=RASPI_AI_AUDIO_LOOP=1

[Install]
WantedBy=multi-user.target
